# price-alert.service
[Unit]
Description=Underbot — WoW price alert service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/spike/github-release-telegramundermine/monitor.py
WorkingDirectory=/home/spike/github-release-telegramundermine
User=spike
EnvironmentFile=/home/spike/github-release-telegramundermine/.env
