# price-alert.timer
[Unit]
Description=Run underbot deal scan every 8 hours

[Timer]
OnCalendar=*-*-* 00,08,16:00:00
Persistent=true

[Install]
WantedBy=timers.target
